# systemd/axel-torrent-cleaner.timer
[Unit]
Description=Start torrent cleanup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
